# configs/automatic-backup/sonarqube-deployment-backup.service
[Unit]
Description=SonarQube Deployment Backup

[Service]
ExecStart=/usr/bin/sonarqube-deployment-backup.sh
